# python/floppyServer.socket
[Unit]
Description=floppyServer socket

[Socket]
ListenStream=80
